# upstart/tracevisor.conf
description	"tracevisor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

console none

script
	cd ~tracevisor
	sudo -u tracevisor /var/lib/tracevisor/tracevisor/tracevisor.py
end script
